# cpu-fan.service
[Unit]
Description=Runs the CPU cooling fan

[Service]
ExecStart=/usr/bin/cpu-fan

[Install]
WantedBy=multi-user.target